# config-file-templates/supervisor/webhook.conf
[program:webhook]
command=webhook -hooks %APP_DIRECTORY%/Misc/hook.json -verbose
stdout_logfile=/var/log/shipper/webhook.log
redirect_stderr=true
